# .claude/voice-server/pai-voice-server.service
[Unit]
Description=PAI Voice Server - Personal AI Voice notification server using ElevenLabs TTS
After=network.target sound.target

[Service]
Type=simple
ExecStart=%h/.bun/bin/bun run %h/.claude/voice-server/server.ts
WorkingDirectory=%h/.claude/voice-server
Restart=on-failure
RestartSec=5
StandardOutput=append:%h/.local/state/pai-voice-server/voice-server.log
StandardError=append:%h/.local/state/pai-voice-server/voice-server-error.log

# Environment
Environment="HOME=%h"
Environment="PATH=%h/.bun/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PORT=8888"

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
